# deploy/campus-celery.service
[Unit]
Description=Campus Celery Worker
After=network.target mysql.service redis.service
Requires=redis.service

[Service]
User=www-data
Group=www-data
WorkingDirectory=/var/www/campus-backend
Environment="PATH=/var/www/campus-backend/venv/bin"
Environment="DJANGO_SETTINGS_MODULE=config.settings"
ExecStart=/var/www/campus-backend/venv/bin/celery \
    -A config worker \
    --loglevel=info \
    --logfile=/var/log/campus/celery.log
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
