# deploy/campus-celery-beat.service
[Unit]
Description=Campus Backend — Celery Beat Scheduler
After=network.target redis.service
Requires=redis.service

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/var/www/campus/campus-backend
Environment=PATH=/var/www/campus/campus-backend/venv/bin
Environment=PYTHONPATH=/var/www/campus/campus-backend
Environment=DJANGO_SETTINGS_MODULE=config.settings
EnvironmentFile=/var/www/campus/campus-backend/.env
ExecStart=/var/www/campus/campus-backend/venv/bin/celery \
    -A config beat \
    --loglevel=info \
    --logfile=/var/log/campus-celery-beat.log \
    --scheduler django_celery_beat.schedulers:DatabaseScheduler
ExecReload=/bin/kill -s HUP $MAINPID
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
